# var/lib/dpkg/info/mxs824-system-config.list
/.
/etc
/etc/apt
/etc/apt/listchanges.conf
/etc/fstab
/etc/fw_env.config
/etc/nginx
/etc/nginx/conf.d
/etc/nginx/conf.d/mxs824.conf
/etc/opt
/etc/opt/ssh
/etc/opt/ssh/sshd_config
/etc/systemd
/etc/systemd/system
/etc/systemd/system/nginx.service.d
/etc/systemd/system/nginx.service.d/nginx-override.conf
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/92-mxs824.rules
/opt
/opt/BIST
/opt/BIST/common.sh
/opt/BIST/i2c-devices.sh
/opt/BIST/local-hw.sh
/opt/BIST/node-errors.js
/opt/BIST/stm-handshake.js
/opt/BIST/system-devices.sh
/opt/BIST/test.sh
/opt/apt
/opt/apt/change_repo.sh
/opt/apt/devel.list
/opt/apt/internal_master_repo.list
/opt/apt/internal_staging_repo.list
/opt/apt/internal_ubuntu_repo.list
/opt/apt/local_directory.list
/opt/apt/public_repo.list
/opt/ltc3887
/opt/ltc3887/ltc3887
/opt/ltc3887/mxs824_ltc.txt
/opt/tps544c20
/opt/tps544c20/mxs824_tps.csv
/opt/tps544c20/tps544c20
/opt/usb_handler.sh
/root
/root/.ssh
/root/.ssh/authorized_keys
/usr
/usr/local
/usr/local/include
/usr/local/include/i2c.h
/usr/local/include/misc.h
/usr/local/include/pfx_crc.h
/usr/local/include/pmc.h
/usr/local/lib
/usr/local/lib/libmxs824.so
/usr/share
/usr/share/doc
/usr/share/doc/mxs824-system-config
/usr/share/doc/mxs824-system-config/changelog.gz
