# etc/systemd/system/nginx.service.d/nginx-override.conf
[Service]
TimeoutStopSec=90
Restart=always
RestartSec=500ms
